# deploy/gunicorn.service.example
# /etc/systemd/system/merrymeal-gunicorn.service — sample.
# Install with:
#   sudo cp deploy/gunicorn.service.example /etc/systemd/system/merrymeal-gunicorn.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now merrymeal-gunicorn

[Unit]
Description=MerryMeal Gunicorn
After=network.target mysql.service

[Service]
Type=notify
User=merrymeal
Group=merrymeal
WorkingDirectory=/srv/merrymeal
EnvironmentFile=/etc/merrymeal/env
ExecStart=/srv/merrymeal/.venv/bin/gunicorn \
    config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
